# home/.config/systemd/user/ollama.service
[Unit]
Description=Ollama Service

[Service]
ExecStart=/usr/sbin/ollama serve

[Install]
WantedBy=default.target
