# home/.config/systemd/user/paccache-cleanup.timer
[Unit]
Description=Run paccache-weekly.service periodically

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
