# home/.config/systemd/user/paccache-cleanup.service
[Unit]
Description=Pacman Cache Cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -ruk1
ExecStart=/usr/bin/paccache -rk2
